# biryaniclub.service
# Production Systemd service file for Biryani Club
# Place this file at: /etc/systemd/system/biryaniclub.service

[Unit]
Description=Biryani Club Web Application
After=network.target

[Service]
Type=simple
User=biryaniclub
Group=biryaniclub
WorkingDirectory=/home/biryaniclub/biryaniclub14
Environment=PATH=/home/biryaniclub/biryaniclub14/venv/bin
EnvironmentFile=/home/biryaniclub/biryaniclub14/.env
ExecStart=/home/biryaniclub/biryaniclub14/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 3 --timeout 120 --max-requests 1000 --preload app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=3

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/biryaniclub/biryaniclub14

[Install]
WantedBy=multi-user.target

